# infra/iac/tofu/templates/iedora.service
[Unit]
Description=iedora compose stack
Requires=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/etc/iedora
ExecStart=/usr/bin/docker compose -f /etc/iedora/docker-compose.yml up -d --remove-orphans
# ExecReload is idempotent — `compose up -d` only recreates services
# whose config changed. Day-2 reconciliation goes through `systemctl
# reload`, NOT `restart`, so unchanged containers (postgres
# mid-transaction, openobserve with live ingest) keep running.
ExecReload=/usr/bin/docker compose -f /etc/iedora/docker-compose.yml up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f /etc/iedora/docker-compose.yml down

[Install]
WantedBy=multi-user.target
